# bin/launchd/com.mirador.ovh-budget.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.mirador.ovh-budget.plist — daily OVH cost alert at 09:00 local.

  Mirror of the GCP budget alert (which uses the native gcloud billing
  budgets mechanism). OVH has no equivalent native API, so this launchd
  job runs `bin/budget/ovh-alert.sh` daily and surfaces a desktop
  notification when month-to-date spend exceeds the €5 threshold.

  Install:
    cp bin/launchd/com.mirador.ovh-budget.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.mirador.ovh-budget.plist

  Remove:
    launchctl unload -w ~/Library/LaunchAgents/com.mirador.ovh-budget.plist
    rm ~/Library/LaunchAgents/com.mirador.ovh-budget.plist

  Inspect:
    launchctl list | grep com.mirador.ovh-budget
    cat /tmp/mirador-ovh-budget.log
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mirador.ovh-budget</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-l</string>
    <string>-c</string>
    <!-- Runs from the svc repo root so REPO_ROOT resolves correctly.
         Sources the user's profile (-l) so OVH_* env vars are picked up. -->
    <string>cd /Users/benoitbesson/dev/workspace-modern/mirador-service &amp;&amp; bin/budget/ovh-alert.sh --quiet</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/mirador-ovh-budget.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/mirador-ovh-budget.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
